# systemd/zerowire-service@.target
[Unit]
Description = ZeroWire service %i online
Requires = zerowire.service
AllowIsolate = no
After = network-online.target
