# systemd/zerowire@.target
[Unit]
Description = ZeroWire peer %i online
Requires = zerowire.service
AllowIsolate = no
After = network-online.target
